# scripts/com.admedical.admin_agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  관리자 대시보드 Agent — 맥북 백그라운드 폴러.

  설치:
    cp scripts/com.admedical.admin_agent.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.admedical.admin_agent.plist

  중지:
    launchctl unload ~/Library/LaunchAgents/com.admedical.admin_agent.plist

  상태:
    launchctl list | grep com.admedical.admin_agent

  로그:
    tail -f /Users/halim/Downloads/아카이브/admedical_website/logs/admin_agent.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.admedical.admin_agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/halim/Downloads/아카이브/admedical_website/venv/bin/python</string>
        <string>/Users/halim/Downloads/아카이브/admedical_website/scripts/admin_agent.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/halim/Downloads/아카이브/admedical_website</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/halim/Downloads/아카이브/admedical_website/logs/admin_agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/halim/Downloads/아카이브/admedical_website/logs/admin_agent.log</string>
</dict>
</plist>
