# scripts/com.admedical.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  매일 새벽 5시(KST) admedical 파이프라인 자동 실행.

  설치 방법 (사장님 1번만 실행):
    cp scripts/com.admedical.daily.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.admedical.daily.plist

  중지 방법:
    launchctl unload ~/Library/LaunchAgents/com.admedical.daily.plist

  상태 확인:
    launchctl list | grep com.admedical.daily

  로그 확인:
    cat /Users/halim/Desktop/admedical_website/logs/daily_*.log

  주의:
    - 맥북이 켜져 있어야 동작 (sleep 모드면 깨우진 않음).
    - 대신 sleep 후에도 5시가 지나서 깨어나면 그 시점에 1번 실행됨.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.admedical.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/halim/Desktop/admedical_website/scripts/daily_pipeline.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/halim/Desktop/admedical_website</string>

    <key>StandardOutPath</key>
    <string>/Users/halim/Desktop/admedical_website/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/halim/Desktop/admedical_website/logs/launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
